# package/usr/lib/systemd/system/openfortivpn-indicator.service
[Unit]
Description=OpenFortiVPN tunnel (controlado pelo indicador da bandeja)
Documentation=man:openfortivpn(1)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Usa exatamente o arquivo de configuração padrão da CLI do openfortivpn.
# stdbuf -oL -eL força line-buffering: sem isso a saída (inclusive a URL
# de login SAML "Authenticate at '...'") fica presa em buffer e não chega
# ao journald a tempo do indicador monitorar e abrir o navegador.
ExecStart=/usr/bin/stdbuf -oL -eL /usr/bin/openfortivpn -c /etc/openfortivpn/config
# Deixa o openfortivpn restaurar as rotas/DNS ao encerrar.
KillSignal=SIGINT
Restart=no

[Install]
WantedBy=multi-user.target
